# local/service/cleanup.timer
[Unit]
Description=Reset Dashboard

[Timer]
OnCalendar=*-*-* 08:10:00

[Install]
WantedBy=timers.target

