# local/service/cleanup.service
[Unit]
Description=Reset Machine

[Service]
Type=simple
User=kermit
ExecStart=reboot
EnvironmentFile=/etc/environment